# infra/egress/hh-egress.service
[Unit]
Description=hh-assistant container egress allowlist (nftables, rendered from src/ops/egress-allowlist.ts)
Documentation=file:///home/hh/hh-assistant/infra/runtime.md
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Re-derive the egress bridge (its name tracks the docker network id) and
# apply the default-deny ruleset (creates the table). BOOT-ONLY now (V2_NOTES
# §11): the periodic re-resolve is hh-egress-refresh.service on the timer, which
# swaps set elements without the fail-open window `apply`'s table-delete opens.
# Idempotent, so re-running by hand is safe.
ExecStart=/usr/bin/env bash -c 'EG=br-$(docker network inspect hh-assistant_egress -f "{{slice .Id 0 12}}") && EGRESS_IFACE="$EG" bash /home/hh/hh-assistant/infra/egress/nftables.sh apply'

[Install]
WantedBy=multi-user.target
