# infra/egress/hh-egress.timer
[Unit]
Description=refresh hh-assistant egress allowlist (re-resolve CDN IPs before the 1h nft set TTL)

[Timer]
# Triggers the REFRESH unit, not the same-named apply unit (V2_NOTES §11):
# `apply` deletes+rebuilds the table (a ~1-2s fail-OPEN window every tick),
# whereas `refresh` swaps only the nft set elements with the table loaded. Boot-
# time table creation stays with hh-egress.service (WantedBy=multi-user.target).
Unit=hh-egress-refresh.service
# 2min after boot lets docker recreate the egress network first; then every
# 15min — 4x margin under the set's 1h element timeout.
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
